# scheduler/valuation_fmp_collector.service
[Unit]
Description=Valuation FMP Collector (financials/estimates/quote -> fmp_financials/fmp_estimates/fmp_metrics)
After=network.target postgresql.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/asset-cloud/scheduler
ExecStart=/usr/bin/python3 /home/ubuntu/asset-cloud/scheduler/valuation_fmp_collector.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
