# scheduler/valuation_fmp_metrics.timer
# ══════════════════════════════════════════════════════════════════════════
# [실행 시각 상수] 08:50 KST (Asia/Seoul)
#
# valuation_sec_edgar_pipeline.timer(08:30 KST), valuation_fmp_collector.timer(08:40 KST)
# 이후에 실행되도록 10분 간격을 유지한다.
# ══════════════════════════════════════════════════════════════════════════
[Unit]
Description=Timer for Valuation FMP Metrics Engine (08:50 KST daily)

[Timer]
OnCalendar=*-*-* 08:50:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
